# state/sss_home/MiG-SSS/iso-helpers/etc/init.d/EXE
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
    after FE
}

DAEMON=/usr/bin/daemon
EXE_PGIDFILE=/var/run/MiG_EXE.pgid
MASTER_PGIDFILE=/var/run/MiG_master.pgid
EXE_BINARY=/usr/bin/start_exe.sh


start() {
	ebegin "Starting ${SVCNAME}"
	    if test -f ${EXE_PGIDFILE} 
	    then
    	       eerror "${SVCNAME} already started"
    	       return 1
    	    else
    	       ${DAEMON} "${EXE_BINARY}" "${EXE_PGIDFILE}" &
    	    fi
	eend $?  
}

stop() {
        ebegin "Stopping ${SVCNAME}"
	    if ! test -f ${EXE_PGIDFILE}
	    then
	       eerror "${SVCNAME} not started"
	       return 1
	    else
	       # Kill all processes in group,
	       # to make sure alle MiG scripts are killed.
	       kill -9 -`cat ${EXE_PGIDFILE}`
	       rm -f ${EXE_PGIDFILE}
	       if test -f ${MASTER_PGIDFILE}
	       then
	          kill -9 -`cat ${MASTER_PGIDFILE}` 
	       	  rm -f ${MASTER_PGIDFILE}
	       fi
	    fi
	eend $?
}
